# .config/autostart/Conky.desktop
[Desktop Entry]
Exec=~/.config/conky/conky.sh
Name=Conky
OnlyShowIn=LXQt;
Type=Application
Version=1.0
