# .system/arch/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
#ExecStart=-/sbin/agetty -a toor -J %I $TERM
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin toor - $TERM
Type=simple
TTYVTDisallocate=no
